# systemd/user/autosync.service
[Unit]
Description=Auto sync MyFish repo

[Service]
Type=oneshot
ExecStart=/home/commander/MyFish/scripts/autosync.sh
